# docker/supervisor/mongo.conf
[program:mongo]
command=/usr/bin/mongod
redirect_stderr=true